# fs/etc/services.d/http-server/run
#!/usr/bin/with-contenv sh

/usr/bin/http-server /home/ -p 80 -c-1 >> /var/log/http-server.log 2>&1
