# fs/etc/services.d/sshd/run
#!/bin/sh

ssh-keygen -A > /dev/null
/usr/sbin/sshd
